# lessons/system services/systemd/demo.service
[Unit]
After=network-online.target
[Service]
Type=exec
ExecStart=/home/student/demo
[Install]
WantedBy=multi-user.target

